# server/payload/systemd/tls-renew.timer
[Unit]
Description=Daily TLS certificate renewal check

[Timer]
OnCalendar=daily
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
